# le-renew.cron
00 07 1 * * root /root/le-renew/le-renew.sh >> /root/le-renew/le-renew.log 2>&1
